# deploy/polish-open-source-rank-crawl-resume.service
[Unit]
Description=Resume interrupted Polish Open Source crawl jobs
After=polish-open-source-rank.service
Wants=polish-open-source-rank.service
StartLimitIntervalSec=0

[Service]
Type=oneshot
WorkingDirectory=/home/ciembor/polish-open-source-rank
EnvironmentFile=/home/ciembor/polish-open-source-rank/.env.local
Restart=on-failure
RestartSec=60
SuccessExitStatus=75
TimeoutStartSec=infinity
TimeoutStopSec=45
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
CPUWeight=20
IOWeight=20
ExecStartPre=/usr/bin/mkdir -p /home/ciembor/polish-open-source-rank/tmp
ExecStartPre=-/usr/bin/podman rm -f polish-open-source-rank-crawl-resume
ExecStart=/usr/bin/flock -n -E 75 /home/ciembor/polish-open-source-rank/tmp/crawl.lock /usr/bin/podman run --name=polish-open-source-rank-crawl-resume --user=1000:1000 --read-only --tmpfs /app/tmp:rw,noexec,nosuid,nodev,size=64m --memory=768m --memory-swap=768m --memory-reservation=512m --cpus=0.5 --cpu-shares=128 --pids-limit=256 --env-file /home/ciembor/polish-open-source-rank/.env.local -e RACK_ENV=production -e APP_BASE_PATH=/ -e BASE_URL=https://polish-open-source.pl -v /home/ciembor/polish-open-source-rank/db:/app/db:rw -v /home/ciembor/polish-open-source-rank/log:/app/log:rw localhost/polish-open-source-rank:latest bundle exec ruby bin/resume_crawls
ExecStop=-/usr/bin/podman stop --time 30 polish-open-source-rank-crawl-resume
ExecStopPost=-/usr/bin/podman rm -f polish-open-source-rank-crawl-resume

[Install]
WantedBy=multi-user.target
