# deploy/polish-open-source-rank-alerts.service
[Unit]
Description=Polish Open Source Rank production alert monitor
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/ciembor/polish-open-source-rank
EnvironmentFile=-/home/ciembor/polish-open-source-rank/.env.local
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ciembor/polish-open-source-rank/tmp /home/ciembor/polish-open-source-rank/log
RestrictSUIDSGID=true
LockPersonality=true
ExecStart=/usr/bin/python3 /home/ciembor/polish-open-source-rank/scripts/production_alert_monitor.py
