# deploy/polish-open-source-rank-monthly.service
[Unit]
Description=Polish Open Source Rank monthly snapshot
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/ciembor/polish-open-source-rank
EnvironmentFile=/home/ciembor/polish-open-source-rank/.env.local
Restart=on-failure
RestartSec=5
ExecStart=/usr/bin/podman run --rm --env-file /home/ciembor/polish-open-source-rank/.env.local -e RACK_ENV=production -e APP_BASE_PATH=/ -e BASE_URL=https://polish-open-source.pl -v /home/ciembor/polish-open-source-rank/db:/app/db -v /home/ciembor/polish-open-source-rank/log:/app/log localhost/polish-open-source-rank:latest bundle exec ruby bin/monthly_rankings
