# deploy/polish-github-rank-monthly.service
[Unit]
Description=Polish GitHub Rank monthly snapshot
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/ciembor/polish-github-rank
EnvironmentFile=/home/ciembor/polish-github-rank/.env.local
ExecStart=/usr/bin/podman run --rm --env-file /home/ciembor/polish-github-rank/.env.local -e RACK_ENV=production -e APP_BASE_PATH=/polish-github-rank -e BASE_URL=https://maciej.ciemborowicz.eu/polish-github-rank -v /home/ciembor/polish-github-rank/db:/app/db -v /home/ciembor/polish-github-rank/log:/app/log localhost/polish-github-rank:latest bundle exec ruby bin/monthly_rankings
