# deploy/polish-open-source-rank.service
[Unit]
Description=Polish Open Source Rank web app
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ciembor/polish-open-source-rank
EnvironmentFile=/home/ciembor/polish-open-source-rank/.env.local
Restart=always
RestartSec=5
TimeoutStartSec=90
CPUWeight=200
IOWeight=200
OOMScoreAdjust=-500

ExecStartPre=-/usr/bin/podman rm -f polish-open-source-rank
ExecStart=/usr/bin/podman run --name=polish-open-source-rank --user=1000:1000 --read-only --tmpfs /app/tmp:rw,noexec,nosuid,nodev,size=64m --memory=1200m --memory-swap=1200m --memory-reservation=512m --cpus=1.0 --pids-limit=256 -p 127.0.0.1:9293:9293 --env-file /home/ciembor/polish-open-source-rank/.env.local -e RACK_ENV=production -e APP_BASE_PATH=/ -e BASE_URL=https://polish-open-source.pl -v /home/ciembor/polish-open-source-rank/db:/app/db:rw -v /home/ciembor/polish-open-source-rank/log:/app/log:rw localhost/polish-open-source-rank:latest
ExecStartPost=-/usr/bin/systemctl start --no-block polish-open-source-rank-crawl-resume.service
ExecStop=/usr/bin/podman stop --time 30 polish-open-source-rank
ExecStopPost=-/usr/bin/podman rm -f polish-open-source-rank
KillMode=mixed
TimeoutStopSec=45

[Install]
WantedBy=multi-user.target
