# deploy/polish-open-source-rank-monitor.timer
[Unit]
Description=Polish Open Source Rank server health monitor

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
