# deploy/polish-github-rank-monthly.timer
[Unit]
Description=Run Polish GitHub Rank monthly snapshot

[Timer]
OnCalendar=*-*-02 03:15:00
Persistent=true

[Install]
WantedBy=timers.target

