# deploy/polish-open-source-rank-crawl.service
[Unit]
Description=Polish Open Source Rank manual crawl
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=0

[Service]
Type=oneshot
WorkingDirectory=/home/ciembor/polish-open-source-rank
EnvironmentFile=/home/ciembor/polish-open-source-rank/.env.local
EnvironmentFile=-/home/ciembor/polish-open-source-rank/.crawl.env
Restart=on-failure
RestartSec=60
TimeoutStartSec=infinity
TimeoutStopSec=45
ExecStartPre=/usr/bin/mkdir -p /home/ciembor/polish-open-source-rank/tmp
ExecStartPre=-/usr/bin/podman rm -f polish-open-source-rank-crawl
ExecStart=/usr/bin/flock -n /home/ciembor/polish-open-source-rank/tmp/crawl.lock /usr/bin/podman run --name=polish-open-source-rank-crawl --memory=1200m --memory-swap=1400m --cpus=1.0 --env-file /home/ciembor/polish-open-source-rank/.env.local -e CRAWL_ARGS -e RACK_ENV=production -e APP_BASE_PATH=/ -e BASE_URL=https://polish-open-source.pl -v /home/ciembor/polish-open-source-rank/db:/app/db -v /home/ciembor/polish-open-source-rank/log:/app/log localhost/polish-open-source-rank:latest sh -lc 'bundle exec ruby bin/monthly_rankings ${CRAWL_ARGS:-}'
ExecStop=-/usr/bin/podman stop --time 30 polish-open-source-rank-crawl
ExecStopPost=-/usr/bin/podman rm -f polish-open-source-rank-crawl
